# Script/mydaemon.service
--/etc/systemd/system/
[Unit]
Description=My_daemon

[Service]
ExecStart=/usr/bin/mydaemon
Restart=on-failure

[Install]
WantedBy=multi-user.target
